# vault.service
[Unit]
Description=vault
Requires=network.target

[Service]
ExecStart=/home/vagrant/vault server -config=/home/vagrant/sync/conf/vault.hcl

[Install]
WantedBy=multi-user.target
